# templates/com.easewise.codex-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.easewise.codex-worker</string>
  <key>UserName</key>
  <string>__USER__</string>
  <key>GroupName</key>
  <string>__GROUP__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__WORKER_BIN__</string>
    <string>--config</string>
    <string>__CONFIG__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__APP_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>CODEX_HOME</key>
    <string>__CODEX_HOME__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__LOG_ROOT__/worker.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_ROOT__/worker.stderr.log</string>
</dict>
</plist>
